# roles/headless_game_streaming/templates/headless-game-streaming.service.j2
[Unit]
Description=Headless game streaming stack (Sunshine + Steam + gamescope)
After=network-online.target systemd-user-sessions.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User={{ headless_game_streaming_user }}
Group={{ headless_game_streaming_group }}
SupplementaryGroups={{ headless_game_streaming_extra_groups | join(' ') }}
WorkingDirectory={{ headless_game_streaming_home }}
{% for key, value in headless_game_streaming_service_environment.items() %}
Environment={{ key }}={{ value }}
{% endfor %}
RuntimeDirectory={{ headless_game_streaming_service_name }} user-{{ headless_game_streaming_uid }}
RuntimeDirectoryMode=0755
StateDirectory={{ headless_game_streaming_service_name }}
StateDirectoryMode=0755
ExecStart=/usr/local/bin/headless-game-streaming-service
ExecStop=/bin/sh -lc '/usr/bin/pkill -KILL -g $MAINPID || true; /usr/bin/pkill -TERM -u {{ headless_game_streaming_user }} sunshine || true; /usr/bin/pkill -TERM -u {{ headless_game_streaming_user }} gamescope || true; /usr/bin/pkill -TERM -u {{ headless_game_streaming_user }} steam || true'
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=15
Restart=on-failure
RestartSec=2
NoNewPrivileges=false
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
