# roles/dashboard/templates/homer.service.j2
[Unit]
Description=Systemd service to start a Homer dashboard at boot.

[Service]
Type=simple
ExecStart=/bin/python3 /opt/homer/server.py

[Install]
WantedBy=multi-user.target