# roles/litellm/templates/headroom.service.j2
[Unit]
Description=Headroom Proxy Service
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/opt/litellm/venv/bin/headroom proxy \
    --host 127.0.0.1 \
    --port {{ litellm_headroom_port }}
{%- if litellm_headroom_logging %} \
    --log-file /var/log/litellm/headroom-requests.jsonl
{%- endif %}

User={{ tesseract_username }}
Group={{ tesseract_username }}
WorkingDirectory=/opt/litellm

Restart=always
RestartSec=3

Environment="PATH=/opt/litellm/venv/bin:{{ ansible_facts['env'].PATH }}"
Environment="HEADROOM_TELEMETRY=off"
Environment="HEADROOM_LOG_LEVEL={{ litellm_headroom_log_level }}"

{% if litellm_headroom_logging %}
StandardOutput=append:/var/log/litellm/headroom.log
StandardError=append:/var/log/litellm/headroom.log
{% endif %}

[Install]
WantedBy=multi-user.target
